# supervisor.ini
[inet_http_server]
port = :9001

[program:rig]
command = echo %(8000+process_num)d
autostart = true
numprocs = 12