# supervisor.ini
[inet_http_server]
port = :9001

[program:SignallingWebServer]
numprocs = 8
process_name = %(program_name)s_%(process_num)d
command = ./CityWar3.exe -AudioMixer -PixelStreamingPort=%(8100+process_num)d -Log -RenderOffScreen -ForceRes -ResX=1920 -ResY=1080 -GraphicsAdapter=%(process_num%4)d
priority = 999 # 命令启动的优先级
autorestart = true # 失败后自动启动
directory = ./SignallingWebServer