# supervisor.ini
[inet_http_server]
port = :9001

[program:SignallingWebServer]
numprocs = 1
process_name = %(program_name)s_%(process_num)d
command = serve
priority = 999 # 命令启动的优先级
autorestart = true # 失败后自动启动
cron = 0/5 * * * * ?
autostart = true