# supervisor.ini
[inet_http_server]
port = :9001

[program:rig]
process_name=%(program_name+"_"+process_num)s
command = echo %(program_name)s
autostart = true
numprocs = 1100
